# motioneye-container.service
# container-2193796636fdba68b7c869be05ae35bb3bb47f92ac2b1b2193424aa9f705dd23.service
# autogenerated by Podman 2.0.4
# Mon Aug 10 09:17:23 EDT 2020

[Unit]
Description=Podman container-2193796636fdba68b7c869be05ae35bb3bb47f92ac2b1b2193424aa9f705dd23.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
ExecStart=/usr/bin/podman start 2193796636fdba68b7c869be05ae35bb3bb47f92ac2b1b2193424aa9f705dd23
ExecStop=/usr/bin/podman stop -t 10 2193796636fdba68b7c869be05ae35bb3bb47f92ac2b1b2193424aa9f705dd23
ExecStopPost=/usr/bin/podman stop -t 10 2193796636fdba68b7c869be05ae35bb3bb47f92ac2b1b2193424aa9f705dd23
PIDFile=/var/run/containers/storage/overlay-containers/2193796636fdba68b7c869be05ae35bb3bb47f92ac2b1b2193424aa9f705dd23/userdata/conmon.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
